# python_server/vr-server.service

[Unit]
Description=VR Command Center Server
After=network.target

[Service]
User=vr-kiosk
Group=vr-kiosk
WorkingDirectory=/opt/vr-command-center
ExecStart=/opt/vr-command-center/venv/bin/python server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Environment variables
Environment=VR_SERVER_PORT=8081
Environment=VR_GAMES_CONFIG=games.json
Environment=VR_STATUS_INTERVAL=5
Environment=LOG_LEVEL=INFO
Environment=VR_DEFAULT_SESSION_DURATION=600
Environment=VR_CPU_WARNING_THRESHOLD=80
Environment=VR_MEMORY_WARNING_THRESHOLD=80

# Uncomment and edit if you need to run as a different user
#User=vr-admin
#Group=vr-admin

[Install]
WantedBy=multi-user.target
