# python_server/vr-server.service

[Unit]
Description=VR Command Center Server
After=network.target

[Service]
User=vr-kiosk
Group=vr-kiosk
WorkingDirectory=/opt/vr-command-center
ExecStart=/opt/vr-command-center/venv/bin/python server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Environment variables can be set here
Environment=VR_SERVER_PORT=8081
Environment=VR_GAMES_CONFIG=games.json

[Install]
WantedBy=multi-user.target
